# root/etc/services.d/subdlsrv/run
#!/usr/bin/with-contenv bash

umask 022

cd /app

exec \
    s6-setuidgid abc doppplerr \
        --no-color \
        --verbose \
        --basedir / \
        --port 8086 \
        --appdir /app \
